# launchd/com.ganeshtalluri.github-profile-token-counter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ganeshtalluri.github-profile-token-counter</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/ganeshtalluri/.local/share/github-profile-token-counter/g4nesh/scripts/update-codex-token-counter.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/ganeshtalluri/.local/share/github-profile-token-counter/g4nesh</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>TZ</key>
    <string>America/Phoenix</string>
    <key>TOKEN_COUNTER_START_DATE</key>
    <string>2026-01-01</string>
    <key>CODEX_USAGE_SPEED</key>
    <string>fast</string>
    <key>CCUSAGE_VERSION</key>
    <string>20.0.14</string>
    <key>GIT_TERMINAL_PROMPT</key>
    <string>0</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>3600</integer>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>/Users/ganeshtalluri/.local/state/github-profile-token-counter/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/ganeshtalluri/.local/state/github-profile-token-counter/launchd.err.log</string>
</dict>
</plist>
